# update_uptime.cron
*/10 * * * * root /usr/bin/update_uptime.sh
