# update_uptime.cron
*/10 * * * * nobody /usr/bin/update_uptime.sh
